# launchd/Library/LaunchAgents/com.aviral.dotfiles-auto-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aviral.dotfiles-auto-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export REPO_DIR="$HOME/dotfiles" TARGET_BRANCH=main SYNC_INTERVAL_SECONDS=300 STATE_DIR="$HOME/.local/state/auto-git-sync/dotfiles"; exec "$HOME/dotfiles/scripts/auto-git-sync" --daemon</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
