# systemd/.config/systemd/user/dotfiles-auto-sync.service
[Unit]
Description=Auto-sync dotfiles with remote every 5 minutes
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dotfiles
Environment=REPO_DIR=%h/dotfiles
Environment=TARGET_BRANCH=main
Environment=SYNC_INTERVAL_SECONDS=300
Environment=STATE_DIR=%h/.local/state/auto-git-sync/dotfiles
ExecStart=%h/dotfiles/scripts/auto-git-sync --daemon
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
