# systemd/.config/systemd/user/pi-telegram.service
[Unit]
Description=Pi over Telegram daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/vault
ExecStart=%h/bin/pi-telegram-daemon.py
Restart=always
RestartSec=5
TimeoutStopSec=10
Environment=PATH=%h/.nvm/versions/node/v22.22.3/bin:%h/.local/bin:%h/go/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%h/.config/pi-telegram.env

[Install]
WantedBy=default.target
